# sandbox/rootfs/etc/supervisor/conf.d/ssh.conf
[program:sshd]
command=/usr/sbin/sshd -D
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd.err
